# rules/persistence/registry_run_key.yml
title: Registry Run Key Persistence
id: b5c3d0e2-7f9a-4b4c-0d3e-2f6a7b8c9d0e
status: experimental
description: |
  Detects modifications to registry Run keys, which are commonly abused by malware
  to achieve persistence across system reboots.
references:
    - https://attack.mitre.org/techniques/T1547/001/
    - https://docs.microsoft.com/en-us/windows/win32/setupapi/run-and-runonce-registry-keys
author: Your Name
date: 2026/01/16
tags:
    - attack.persistence
    - attack.t1547.001
logsource:
    category: registry_event
    product: windows
detection:
    selection:
        TargetObject|contains:
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run'
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce'
            - '\SOFTWARE\WOW6432Node\Microsoft\Windows\CurrentVersion\Run'
            - '\SOFTWARE\WOW6432Node\Microsoft\Windows\CurrentVersion\RunOnce'
    filter_legitimate:
        Details|contains:
            - 'C:\Program Files\'
            - 'C:\Program Files (x86)\'
            - 'C:\Windows\System32\'
    condition: selection and not filter_legitimate
falsepositives:
    - Legitimate software installation
    - System administrators adding startup programs
level: medium
